# packaging/scripts/before-remove
#!/bin/sh

# Stop and disable the service.
systemctl stop owntracks-publisher
systemctl disable owntracks-publisher
