# packaging/scripts/before-remove
#!/bin/sh

# Stop and disable the service.
systemctl stop owntracks-cli-publisher
systemctl disable owntracks-cli-publisher
